# scripts/launchd/com.pitchai.codex.morning.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pitchai.codex.morning</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-c</string>
        <string>cd /Users/sethvanderbijl/codex && PROMPT_FILE="/Users/sethvanderbijl/codex/scripts/prompts/morning_healthcheck_prompt.txt" /Users/sethvanderbijl/codex/scripts/run_codex_prompt.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>5</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/sethvanderbijl/Library/Logs/codex-automation/launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/sethvanderbijl/Library/Logs/codex-automation/launchd.log</string>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
